# root/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc python3 /app/sabnzbd/SABnzbd.py \
    --config-file /config --server 0.0.0.0:8080


