# root/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc python /app/sabnzbd/SABnzbd.py -f /config/sabnzbd.ini -b 0 -s 0.0.0.0:8080
