# root/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

FAMILY=::

if [ "$HAS_IPV6" = "false" ]; then
    FAMILY=0.0.0.0
fi

exec \
    s6-setuidgid abc python3 /app/sabnzbd/SABnzbd.py \
    --config-file /config --server "$FAMILY"
