# services/sabnzbd/run
#!/bin/bash

exec /sbin/setuser abc /usr/bin/sabnzbdplus --config-file /config --server 0.0.0.0:8080
