# root/etc/s6-overlay/s6-rc.d/init-sabnzbd-config/run
#!/usr/bin/with-contenv bash

# permissions
chown abc:abc \
    /config
