# root/etc/s6-overlay/s6-rc.d/init-sabnzbd-config/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /downloads \
    /incomplete-downloads

# permissions
chown abc:abc \
    /downloads \
    /incomplete-downloads

chown -R abc:abc \
    /config
